# share/systemd/test_server.service.erb
[Unit]
Description=<%= lookup('DESCRIPTION') %>

[Service]
User=<%= lookup('SERVICE_USER') %>
Group=<%= lookup('SERVICE_USER') %>
ExecStart=<%= lookup('COMMAND') %>

[Install]
WantedBy=multi-user.target
