# share/systemd/test_server.service
[Unit]
Description=Test server
#Requires=test_server.socket

[Service]
ExecStart=/usr/bin/test_server serve

[Install]
WantedBy=multi-user.target
